# 99-ic7300-clock-sync.rules
# Icom IC-7300 — auto clock sync on USB connect
#
# Matches the Silicon Labs CP210x USB-serial bridge used by the IC-7300.
# Creates a stable /dev/ic7300 symlink and integrates with systemd so the
# sync service starts both on hot-plug AND when already connected at boot.
#
# Using TAG+="systemd" + ENV{SYSTEMD_WANTS} is the correct approach:
#   - RUN+=systemctl races against systemd startup during boot (fails silently)
#   - SYSTEMD_WANTS lets PID 1 start the service once it is ready, regardless
#     of whether the device was present before or after systemd came up.
#
# Install:
#   sudo cp 99-ic7300-clock-sync.rules /etc/udev/rules.d/
#   sudo udevadm control --reload-rules && sudo udevadm trigger --action=add --subsystem-match=tty

ACTION=="add", \
    SUBSYSTEM=="tty", \
    ATTRS{idVendor}=="10c4", \
    ATTRS{idProduct}=="ea60", \
    GROUP="plugdev", \
    MODE="0660", \
    SYMLINK+="ic7300", \
    TAG+="systemd", \
    ENV{SYSTEMD_WANTS}="ic7300-clock-sync.service"
